# fish-osc.service
[Unit]
Description=Fish Motor OSC Controller Service
Documentation=https://github.com/JonnyCo/PeixeCidade
After=network-online.target graphical-session.target
Wants=network-online.target graphical-session.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/PeixeCidade
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

# Wait for display to be ready, then execute startup script
ExecStart=/bin/bash -c 'while ! pgrep -x "Xorg" > /dev/null; do sleep 1; done; exec /home/pi/PeixeCidade/fish-osc-startup.sh'

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging configuration
StandardOutput=append:/home/pi/PeixeCidade/fish-osc-service.log
StandardError=append:/home/pi/PeixeCidade/fish-osc-error.log

# Security settings (relaxed for GUI access)
PrivateTmp=true
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/home/pi/PeixeCidade

[Install]
WantedBy=graphical.target