# fish-osc.service
[Unit]
Description=Fish Motor OSC Controller Service
Documentation=https://github.com/JonnyCo/PeixeCidade
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/PeixeCidade
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=PYTHONUNBUFFERED=1

# Execute startup script with proper environment
ExecStart=/bin/bash -c 'source /home/pi/.profile; cd /home/pi/PeixeCidade; ./fish-osc-startup.sh'

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging configuration
StandardOutput=append:/home/pi/PeixeCidade/fish-osc-service.log
StandardError=append:/home/pi/PeixeCidade/fish-osc-error.log

# Security settings (relaxed for GUI access)
PrivateTmp=false
ProtectSystem=false
ProtectHome=false

[Install]
WantedBy=multi-user.target